# 33homework/roles/go/files/go.service
[Unit]
Description=go

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/opt/go/main

[Install]
WantedBy=multi-user.target
