# 04homework/crontab
0 * * * *	/home/vagrant/scr/analyzer_log.sh

