# 04homework/crontab
0 * * * *	/home/vagrant/scr/analyzer_log.sh 8 8 2>/dev/null
