# 08homework/01-part/pars_log.timer
[Unit]
Description=Timer for run pars_log.service every 30 sec

[Timer]
OnCalendar=*:*:0/30

[Install]
WantedBy=multi-user.target
